# dbus/blivet.service
[Unit]
Description=Blivet Service

[Service]
Type=dbus
BusName=com.redhat.Blivet0
ExecStart=/usr/libexec/blivetd
